# docker/supervisor/queue.conf
[program:queue-worker]
command=php /var/www/html/artisan queue:work --tries=3
user=www
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/queue-worker.log
